# support_files/linux/splintermail.service
[Unit]
Description=splintermail CITM server.

[Service]
Type=simple
ExecStart=QW install_prefix WQ/bin/splintermail citm
StandardOutput=syslog
StandardError=syslog
Restart=on-failure
User=splintermail
Group=splintermail

[Install]
WantedBy=default.target
